# etc/udev/rules.d/99-bigcam-webcam-power.rules
# BigCam: Disable USB autosuspend for video capture devices to prevent frame drops.
# When a V4L2 device appears, set the parent USB device power control to "on".
ACTION=="add", SUBSYSTEM=="video4linux", RUN+="/bin/sh -c 'echo on > /sys/class/video4linux/%k/device/../power/control 2>/dev/null || true'"
